# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nginx]
command=nginx -g "daemon off;"
stdout_logfile=/var/log/nginx/access.log
stderr_logfile=/var/log/nginx/error.log
autorestart=true

[program:email_logger]
command=python3 /app/email_logger.py
directory=/app
stdout_logfile=/var/log/email_logger.log
stderr_logfile=/var/log/email_logger.log
autorestart=true
